# crontab
* 1 * * * python cron.py
* * 1 * * python soc_update.py
